# templates/traefik2.service.epp
[Unit]
Description=traefik2
Documentation=https://docs.traefik.io/
Wants=basic.target
After=basic.target network.target

[Service]
Type=simple
WorkingDirectory=<%= $traefik2::data_dir %>
ExecStart=<%= $traefik2::bin_dir %>/traefik --configFile=<%= $traefik2::config_dir %>/config.yaml
Restart=on-failure
RestartSec=30s
<% if $traefik2::systemd_workdir != undef { -%>
WorkingDirectory=<%= $traefik2::systemd_workdir %>
<% } -%>

[Install]
WantedBy=multi-user.target
